# autodeployment.service
[Unit]
Description=Gitea Webhook Auto-Deployment Service
After=network.target network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
User=trevor
Group=trevor
WorkingDirectory=/home/trevor/autodeployer
Environment="PATH=/home/trevor/.local/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/trevor/autodeployer/.env
ExecStart=/home/trevor/.local/bin/uv run /home/trevor/autodeployer/webhook_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autodeployment

# Security hardening
NoNewPrivileges=false
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
ReadWritePaths=/home/trevor/.local/log/autodeploymentservice /home/trevor/ManagedFileTransfer /home/trevor/.cache/uv /home/trevor/autodeployer

# Resource limits
MemoryMax=512M
TasksMax=50

[Install]
WantedBy=multi-user.target
