# autodeployment.service
[Unit]
Description=Gitea Webhook Auto-Deployment Service
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trevor
Group=trevor
WorkingDirectory=/home/trevor/autoupdater
Environment="PATH=/home/trevor/.local/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/trevor/autoupdater/.env
ExecStart=/home/trevor/.local/bin/uv run /home/trevor/autoupdater/webhook_server.py
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autodeployment

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
ReadWritePaths=/home/trevor/.local/log/autodeploymentservice /home/trevor/ManagedFileTransfer

# Resource limits
MemoryMax=512M
TasksMax=50

[Install]
WantedBy=multi-user.target
